# profile.d/10-ssh-agent.sh
/home/phaedrus/.profile.d/10-ssh-agent.sh