# profile.d/11-gpg-agent.sh
#!/usr/bin/env sh
# /qompassai/Shell/.profile.d/11-gpg-agent.sh
# Copyright (C) 2025 Qompass AI, All rights reserved
####################################################
export GNUPGHOME="$HOME/.gnupg"
export GPG_TTY=$(tty)
