# profile.d/11-gpg-agent.sh
/home/phaedrus/.profile.d/11-gpg-agent.sh