# profile.d/11-gpg-agent.sh
export GNUPGHOME="$HOME/.gnupg"
export GPG_TTY=$(tty)
