# profile.d/14-cplusinclude.sh
#export CPLUS_INCLUDE_PATH=$CPLUS_INCLUDE_PATH:$OPTIX_INSTALL_DIR/include
