# profile.d/20-hyprland.sh
/home/phaedrus/.profile.d/20-hyprland.sh