# profile.d/21-ue.sh
export UE_LocalDataCachePath="/home/phaedrus/.cache/UnrealEngine/DDC"
export UE_GITDEPS_ARGS="--threads=12 --cache=$HOME/.cache/UnrealEngine/GitDeps --cache-days=180 --cache-size-multiplier=5 --max-retries=10 --all"
