# profile.d/30-hf.sh
export HF_TOKEN=$(pass show hf)
alias hfd='HF_TOKEN=$(pass show hf) hfdownloader'
