# profile.d/41-postgresql.sh
/home/phaedrus/.profile.d/41-postgresql.sh