# profile.d/43-dovecot.sh
export DOVECOT_CONF="$HOME/.config/dovecot/dovecot.conf"
