# profile.d/44-git.sh
#!/usr/bin/env sh
# /qompassai/Shell/.profile.d/44-git.sh
# Copyright (C) 2025 Qompass AI, All rights reserved
####################################################

export GIT_DISCOVERY_ACROSS_FILESYSTEM=1
export GIT_EDITOR='nvim'
export GIT_PAGER='nvimpager'
export GIT_TERMINAL_PROMPT=1
export JJ_EDITOR='nvim'
export JJ_PAGER='nvimpager'
export JJ_TERMINAL_PROMPT=1

