# profile.d/55-lua.sh
export LUA_VERSION=luajit
export LUA_PATH="/usr/local/share/luajit-2.1.0-beta3/?.lua;/usr/local/share/luajit-2.1.0-beta3/?/init.lua;;"
export LUA_CPATH="/usr/local/lib/?.so;;"
