# profile.d/56-python.sh
export PYTHON=python3
export WORKON_HOME="$HOME/.virtualenvs"
export PYTHONSTARTUP="$HOME/.config/pythonrc.py"
