# profile.d/62-bash-comp.sh
/home/phaedrus/.profile.d/62-bash-comp.sh