# profile.d/64-nvim.sh
#!/usr/bin/env sh
# /qompassai/Shell/.profile.d/64-nvim.sh
# ~/.profile.d/64-nvim.sh
# Copyright (C) 2025 Qompass AI, All rights reserved
# -----------------------------

export BAT_THEME="ansi"
export EDITOR="nvim"
export GIT_EDITOR="nvim"
export LANG="en_US.UTF-8"
export LC_ALL="en_US.UTF-8"
export MANPAGER="nvimpager"
export NVIM_LOG_FILE="$XDG_CACHE_HOME/nvim/log"
export NVIMPAGER_CONFIG="$HOME/.config/nvimpager/init.lua"
export NVIMPAGER_FZF="fzf"
export NVIMPAGER_LOG="$HOME/.cache/nvimpager.log"
export NVIMPAGER_NVIM="nvim"
export PAGER="nvimpager"
export PATH="$HOME/.local/share/nvim/mason/bin:$PATH"
export VISUAL="nvim"
export XDG_CACHE_HOME="${XDG_CACHE_HOME:-$HOME/.cache}"
export XDG_CONFIG_HOME="${XDG_CONFIG_HOME:-$HOME/.config}"
export XDG_DATA_HOME="${XDG_DATA_HOME:-$HOME/.local/share}"

