# profile.d/70-terminal.sh
/home/phaedrus/.profile.d/70-terminal.sh