# profile.d/71-radicle.sh
export PATH="$PATH:/home/$USER/.radicle/bin"
